# linux/dump_timer.timer
[Unit]
Description=Dump Backup Timer

[Timer]
Unit=dump_service.service
OnUnitActiveSec=1s

[Install]
WantedBy=timers.target
